# .config/udev/70-persistent-usb.rules
SUBSYSTEM=="power_supply", ACTION=="change", ENV{POWER_SUPPLY_TYPE}=="Mains", RUN+="/usr/bin/pkill -SIGRTMIN+1 sway-status"
